# roles/desktop/files/config/conky/conky.service
[Unit]
Description=conky daemon

[Service]
ExecStart=conky

[Install]
WantedBy=default.target
